# docker-cleanup.timer
[Unit]
Description=Docker cleanup timer

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
